# systemd/.config/systemd/user/mail.service
[Unit]
Description=Run mail-notify script

[Service]
Type=oneshot
ExecStart=%h/.scripts/mail-notify.sh
