# systemd/.config/systemd/user/sundial.service
[Unit]
Description=Run sundial script

[Service]
Type=oneshot
ExecStart=%h/.scripts/sundial
